# deploy/roger-deploy.timer
[Unit]
Description=Poll GHCR for new roger images and redeploy

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
